# etc/cron.d/namecoinToBind
# /etc/cron.d/namecoinToBind: crontab for namecoinToBin

# Check for new names each 2 minutes
*/2 *		* * *	www-data	/usr/bin/php /home/NamecoinToBind/namescan.php > /dev/null

